# script/supervisor/inowas.interpolation.rpc.server.conf
[program:inowas.interpolation.rpc.server]
command=python /var/www/pyModelling/inowas.interpolation.rpc.server.py
